# docker/supervisor_confs/route_server.ini
[program:route_server]
directory=/opt/tt_source_code/server/im-server-test/%(program_name)s
command=/opt/tt_source_code/server/im-server-test/%(program_name)s/%(program_name)s
autostart=true
autorestart=true
environment=LD_LIBRARY_PATH=/opt/tt_source_code/server/im-server-test/lib